# deploy-agent.service
[Unit]
Description=Deploy Agent HTTP Service
After=network.target

[Service]
Type=simple
User=openclaw
WorkingDirectory=/home/openclaw/apps/deploy-agent
ExecStart=/usr/bin/python3 -m uvicorn main:app --host 0.0.0.0 --port 18795
EnvironmentFile=/home/openclaw/apps/deploy-agent/.env
Restart=always
RestartSec=5

[Install]
WantedBy=multi-user.target
